# etc/supervisor/conf.d/sshd.conf
[supervisord]
nodaemon=true
[program:vncd_go_go]
command=bash -c '/dockerstartup/go_go.sh'
directory=/root/
